# roles/packet_steering/templates/packet-steering-init.j2
#!/bin/sh /etc/rc.common
# /etc/init.d/packet-steering
# Generated by Ansible - Do not edit manually
#
# Custom IRQ affinity and packet steering for NanoPi R6S
# Optimized for CAKE SQM at 1 Gbps symmetric
#
# Based on:
# - https://github.com/StarWhiz/NanoPi-R6S-CPU-Optimization-for-Gigabit-SQM
# - https://forum.openwrt.org/t/sqm-optimization-for-nanopi-r6s-with-24-10-4/243769
#
# CPU Topology (RK3588S):
#   CPU 0-3: Cortex-A55 (efficiency cores, 1.8GHz)
#   CPU 4-7: Cortex-A76 (performance cores, 2.4GHz)
#
# CPU affinity masks:
#   0x01 = CPU 0 (A55)     0x10 = CPU 4 (A76)
#   0x02 = CPU 1 (A55)     0x20 = CPU 5 (A76)
#   0x04 = CPU 2 (A55)     0x40 = CPU 6 (A76)
#   0x08 = CPU 3 (A55)     0x80 = CPU 7 (A76)
#   0x0f = CPU 0-3 (all A55)
#   0xf0 = CPU 4-7 (all A76)

START=99

# Set IRQ affinity for an interface
set_irq_affinity() {
    local iface="$1"
    local affinity="$2"

    for irq in $(grep "${iface}" /proc/interrupts | awk -F: '{print $1}' | tr -d ' '); do
        if [ -n "$irq" ] && [ -f "/proc/irq/$irq/smp_affinity" ]; then
            echo "$affinity" > /proc/irq/$irq/smp_affinity 2>/dev/null
        fi
    done
}

# Set RPS (Receive Packet Steering) for an interface
set_rps() {
    local iface="$1"
    local cpus="$2"

    for queue in /sys/class/net/$iface/queues/rx-*/rps_cpus; do
        [ -f "$queue" ] && echo "$cpus" > "$queue" 2>/dev/null
    done

    # Set RFS flow count if supported
    for queue in /sys/class/net/$iface/queues/rx-*/rps_flow_cnt; do
        [ -f "$queue" ] && echo 4096 > "$queue" 2>/dev/null
    done
}

start() {
    # Enable global RFS (Receive Flow Steering)
    if [ -f /proc/sys/net/core/rps_sock_flow_entries ]; then
        echo 32768 > /proc/sys/net/core/rps_sock_flow_entries
    fi

    # =========================================================================
    # WAN interface ({{ network_interfaces.wan }}) - 2.5 Gbps primary traffic path
    # =========================================================================
    # IRQ handling: CPU 4 (0x10)
    # RPS softirq:  CPU 6 (0x40)
    set_irq_affinity "{{ network_interfaces.wan }}" "10"
    set_rps "{{ network_interfaces.wan }}" "40"

    # =========================================================================
    # LAN trunk interface ({{ network_interfaces.lan_trunk }}) - 2.5 Gbps to switch
    # =========================================================================
    # IRQ handling: CPU 5 (0x20)
    # RPS softirq:  CPU 7 (0x80)
    set_irq_affinity "{{ network_interfaces.lan_trunk }}" "20"
    set_rps "{{ network_interfaces.lan_trunk }}" "80"

    # =========================================================================
    # Management interface ({{ network_interfaces.mgmt }}) - 1 Gbps low priority
    # =========================================================================
    # IRQ handling: CPU 0 (0x01)
    # RPS softirq:  CPUs 1-3 (0x0e)
    set_irq_affinity "{{ network_interfaces.mgmt }}" "01"
    set_rps "{{ network_interfaces.mgmt }}" "0e"

    logger -t packet-steering "Configuration applied"
}

reload_service() {
    start
}
