# roles/packet_steering/templates/packet_steering-disabled.j2
#!/bin/sh /etc/rc.common
# /etc/init.d/packet_steering
# Generated by Ansible - Do not edit manually
#
# This is a no-op replacement for OpenWrt's default packet_steering.
# Custom packet steering is handled by /etc/init.d/packet-steering instead.
#
# The default OpenWrt packet_steering would override our IRQ/RPS settings
# via procd triggers and /usr/libexec/network/packet-steering.uc

START=25
USE_PROCD=1

start_service() {
    return 0
}

reload_service() {
    return 0
}

service_triggers() {
    return 0
}
