# roles/network/templates/trunk-carrier-detect.j2
#!/bin/sh /etc/rc.common
# Bring up trunk interface so the kernel detects carrier state and the LED correctly
# reflects whether a cable is connected. Without this, the bridge keeps the interface
# UP regardless of carrier, and the netdev LED trigger never shows NO-CARRIER.

START=99

boot() {
	ip link set {{ network_interfaces.lan_trunk }} up 2>/dev/null
}

start() {
	boot
}

stop() {
	return 0
}
